# tykj-operation/tykj-operation/MarketSearchDownloader/scripts/downloader-service.conf
# MarketSearchDownloader - service
description     "MarketSearchDownloader - service"
version         "1.0"
author          "gmliao"

# running level.
start on runlevel [2345]
stop on runlevel [06]

respawn


script
    RUNNER="/var/app/enabled/MarketSearchDownloader/scripts/start-downloader-service.sh"
    OPTS="--no-console"
    exec "${RUNNER}" ${OPTS}
end script
